# launchd/com.reelfin.reeltranscode.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.reelfin.reeltranscode.watch</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/uv</string>
      <string>run</string>
      <string>reeltranscode</string>
      <string>--config</string>
      <string>/Users/flo/Documents/GitHub/ReelTranscode/config/reeltranscode.giant_boy.yaml</string>
      <string>watch</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/flo/Documents/GitHub/ReelTranscode</string>

    <key>StandardOutPath</key>
    <string>/Users/flo/Documents/GitHub/ReelTranscode/reports/giant_boy/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/flo/Documents/GitHub/ReelTranscode/reports/giant_boy/launchd.stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
